# leopardo-rh-docs/04_CICD_ET_CONFIG/leopardo-horizon.supervisor.conf
# /etc/supervisor/conf.d/leopardo-horizon.conf
# Configuration Supervisor pour Laravel Horizon (Queue Workers)
# Copier ce fichier sur le serveur o2switch après déploiement initial

[program:leopardo-horizon]
process_name=%(program_name)s
command=php /var/www/leopardo-rh-api/artisan horizon
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/leopardo-horizon.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=3600
startsecs=5
startretries=3

# Après modification de ce fichier :
# sudo supervisorctl reread
# sudo supervisorctl update
# sudo supervisorctl start leopardo-horizon:*
# sudo supervisorctl status
